# root/etc/systemd/system/cockpit.service.d/nethserver.conf
[Service]
EnvironmentFile=-/etc/sysconfig/cockpit
